# sys/hddtemp-exporter.timer
[Unit]
Description=Update hddtemp metrics

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
